# kernel/linux-2.6.36/usr/ramfs/etc/inittab
ttyS0::askfirst:-/bin/sh --login
::sysinit:/etc/rcS
